# MicroOS/mkosi.microos/mkosi.extra/usr/lib/systemd/system/systemd-growfs-root.service.d/override.conf
[Service]
ExecStart= 
ExecStart=/usr/lib/systemd/systemd-growfs /etc

